# macos/org.unicornscan.ChmodBPF.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  org.unicornscan.ChmodBPF.plist
  LaunchDaemon that sets /dev/bpf* permissions at system startup so that
  members of the 'unicornscan' group can capture packets without root.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.unicornscan.ChmodBPF</string>

    <key>Program</key>
    <string>/usr/local/bin/ChmodBPF</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/unicornscan-chmodbpf.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/unicornscan-chmodbpf.log</string>
</dict>
</plist>
